# build_files/systemd/bootc-fetch-apply-updates.service.d/10-emryk.conf
# SECURITY-TODO #7: drop the stock `--apply` flag so the timer fetches and
# stages updates only — never reboots on its own. The user reboots when
# their workload allows. The next boot picks up whatever has been staged.
# Override the entire ExecStart by clearing it first.
[Service]
ExecStart=
ExecStart=/usr/bin/bootc upgrade --quiet
